# config/supervisord.conf
[unix_http_server]
file=%(here)s/../run/supervisor.sock

[inet_http_server]
port = 0.0.0.0:9001
username = bdwilson
password = bdwilson

[supervisord]
pidfile=%(here)s/../run/supervisord.pid
logfile=%(here)s/../log/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
nodaemon=false
minfds=1024
minprocs=200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix://%(here)s/../run/supervisor.sock

[program:elasticsearch]
directory=/home/bdwilson/elasticsearch/bin
command=/home/bdwilson/elasticsearch/bin/elasticsearch
process_name=%(program_name)s
priority=1
numprocs=1
numprocs_start=0
redirect_stderr=true
stdout_logfile=%(here)s/../log/%(program_name)s.log
startsecs=10

[program:facetsearch-podaac]
directory=/home/bdwilson/ops/facetsearch-podaac
#environment=PODAAC_ENV=dev
#command=/home/bdwilson/ops/facetsearch-podaac/manage.py server -h 0.0.0.0 -p 3001
command=gunicorn -w4 -b 0.0.0.0:8889 -k gevent 'podaac:create_app("podaac.settings.ProdConfig", env="prod")'
process_name=%(program_name)s
priority=2
numprocs=1
numprocs_start=0
redirect_stderr=true
stdout_logfile=%(here)s/../log/%(program_name)s.log
startsecs=10

[program:facetview-memex]
directory=/home/bdwilson/ops/facetview-memex
#environment=MEMEX_ENV=dev
#command=/home/bdwilson/ops/facetview-memex/manage.py server -t 0.0.0.0 -p 8888
command=gunicorn -w4 -b 0.0.0.0:8888 -k gevent 'memex:create_app("memex.settings.ProdConfig", env="prod")'
process_name=%(program_name)s
priority=2
numprocs=1
numprocs_start=0
redirect_stderr=true
stdout_logfile=%(here)s/../log/%(program_name)s.log
startsecs=10
